# admin_panel/celery-beat.service
[Unit]
Description=Celery Beat Scheduler for Django Admin Panel
After=network.target redis.target

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/var/www/mindful-journal/admin_panel

# Виртуальное окружение
Environment="PATH=/var/www/mindful-journal/admin_panel/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=admin_panel.settings"

# Загружаем переменные из .env
EnvironmentFile=/var/www/mindful-journal/admin_panel/.env

# Запуск Celery Beat (планировщик)
ExecStart=/var/www/mindful-journal/admin_panel/venv/bin/celery \
    -A admin_panel beat \
    --loglevel=info \
    --pidfile=/var/run/celery/beat.pid \
    --logfile=/var/log/celery/beat.log \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler \
    --detach

ExecStop=/bin/kill -s TERM $MAINPID

# Перезапуск при падении
Restart=on-failure
RestartSec=10

# Создаём директории
RuntimeDirectory=celery
LogsDirectory=celery

[Install]
WantedBy=multi-user.target
